# odootools/bash/openerp-scripts/test-db/openerp.cron
# Copy [DB1] to [DB2]
*/30 * * * * root /etc/openerp/[BRANCH]/server/[DB1]-to-[DB2].sh
